# spec/resources/export/supervisor/app-alpha.conf
[program:app-alpha]
command=./alpha
directory=/tmp/app
user=app
