# spec/resources/export/supervisor/app-bravo.conf
[program:app-bravo]
command=./bravo
directory=/tmp/app
process_name=%(program_name)-%(process_num)
num_procs=2
user=app
